# dot_config/systemd/user/theme-dark.timer
[Unit]
Description=Switch to dark theme at 8:30pm

[Timer]
OnCalendar=*-*-* 20:30:00
Persistent=true

[Install]
WantedBy=timers.target
