# dot_config/systemd/user/theme-light.timer
[Unit]
Description=Switch to light theme at 6am

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true

[Install]
WantedBy=timers.target
