# dot_config/systemd/user/theme-light.service
[Unit]
Description=Switch to light green theme
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/nyx-theme apply light_green

[Install]
WantedBy=default.target
