# dot_config/systemd/user/theme-dark.service
[Unit]
Description=Switch to dark green theme
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/nyx-theme apply dark_green

[Install]
WantedBy=default.target
