# dot_config/systemd/user/theme-dark.timer
[Unit]
Description=Switch to dark theme at 7:30pm

[Timer]
OnCalendar=*-*-* 17:30:00
Persistent=true

[Install]
WantedBy=timers.target
